# root/etc/services.d/web/finish
#!/usr/bin/with-contenv bash

/usr/sbin/apachectl -k graceful

